# V2ex/crontab.cron
00 11 * * * /user/src/V2ex/start_spider.sh