# V2ex/crontab.cron
0 11 * * * /usr/src/v2ex/V2ex/start_spider.sh

